# systemd/kr-shadow-daily.timer
[Unit]
Description=KR Shadow League daily finalization at 16:00 KST (07:00 UTC)
Requires=kr-shadow-daily.service

[Timer]
# 16:00 KST = 07:00 UTC. Verified with systemd-analyze calendar on AWS.
OnCalendar=Mon..Fri *-*-* 07:00:00 UTC
Persistent=true
Unit=kr-shadow-daily.service

[Install]
WantedBy=timers.target
